# scanner.service
[Unit]
Description=Mom Trader Scanner
After=network.target

[Service]
EnvironmentFile=/opt/scanner/.env
ExecStart=/opt/scanner/scanner_unix %H
WorkingDirectory=/opt/scanner
Restart=on-failure
RestartSec=5
User=root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
